# openclaw-prometheus-exporter/com.openclaw.exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.exporter</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/feynman/repos/openclaw-prometheus-exporter/openclaw_exporter.py</string>
    <string>--port</string>
    <string>9101</string>
    <string>--interval</string>
    <string>300</string>
    <string>--providers</string>
    <string>claude</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/openclaw-exporter.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/openclaw-exporter.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
